# systemd/suture-platform.service
[Unit]
Description=Suture Platform
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=suture
Group=suture
WorkingDirectory=/var/lib/suture
Environment=RUST_LOG=info
Environment=SUTURE_DATA_DIR=/var/lib/suture
EnvironmentFile=-/etc/suture/env
ExecStart=/usr/local/bin/suture-platform --addr 0.0.0.0:8080
Restart=always
RestartSec=5
LimitNOFILE=65536
ProtectSystem=strict
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
